# wolserver.service
[Unit]
Description=Wake-on-Lan Server
After=network-online.target

[Service]
Type=simple
ExecStart=/root/wakeonlanserver -c /root/wol-config.json

[Install]
WantedBy=multi-user.target
